# deploy/wms-worker.service
[Unit]
Description=WMS Kaspi Poller Worker
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=nurkhat
Group=nurkhat
WorkingDirectory=/home/nurkhat/wms-backend
EnvironmentFile=/home/nurkhat/wms-backend/.env
ExecStart=/home/nurkhat/wms-backend/.venv/bin/python -m app.workers.scheduler
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wms-worker

[Install]
WantedBy=multi-user.target
